# projects/data_lake/hadoop_azkaban/scheduling/crontab/crontab
@reboot . $HOME/.bashrc; /opt/hadoop/sbin/start-all.sh
@reboot . $HOME/.bashrc; nohup hiveserver2 &
@reboot . $HOME/.bashrc; $AZKABAN_HOME/bin/start-solo.sh